# agent/browser/rootfs/etc/s6-overlay/s6-rc.d/svc-desktop/run
#!/command/with-contenv bash

export DISPLAY=:0
export HOME=/home/claworc
export LANG=en_US.UTF-8

# Wait for X server to be ready
while ! xdpyinfo -display :0 >/dev/null 2>&1; do sleep 0.5; done

# No window manager: the browser is the only X client. Without a WM,
# --start-maximized is a no-op (it's a WM hint), so we size and position the
# Chromium window explicitly via --window-size / --window-position to match
# the Xtigervnc display geometry.

# Build optional extra flags (array preserves spaces in values)
EXTRA_ARGS=()
if [ -n "$CHROMIUM_USER_AGENT" ]; then
  EXTRA_ARGS+=("--user-agent=$CHROMIUM_USER_AGENT")
fi

# Auto-detect browser
if command -v google-chrome-stable >/dev/null 2>&1; then
  BROWSER=google-chrome-stable
  EXTRA_ARGS+=(--no-sandbox)
elif command -v brave-browser >/dev/null 2>&1; then
  BROWSER=brave-browser
  EXTRA_ARGS+=(--no-sandbox)
else
  BROWSER=/usr/bin/chromium
  EXTRA_ARGS+=(--no-sandbox)
fi

PROFILE_DIR=/home/claworc/chrome-data

# Clear stale Singleton* files left behind when Chromium dies without a clean
# shutdown (e.g. container kill). The lock symlinks `<host>-<pid>` and the
# matching socket would otherwise abort every subsequent launch with
# "profile in use", trapping us in a tight respawn loop.
clear_stale_singletons() {
  if pgrep -u claworc -f /usr/lib/chromium/chromium >/dev/null 2>&1; then
    return
  fi
  if pgrep -u claworc -x chrome >/dev/null 2>&1; then
    return
  fi
  rm -f "$PROFILE_DIR/SingletonLock" "$PROFILE_DIR/SingletonCookie" "$PROFILE_DIR/SingletonSocket"
}

# Run browser in a restart loop
while true; do
  clear_stale_singletons
  # CDP binds to 127.0.0.1 only — Chromium 134+ silently refuses non-loopback
  # remote-debugging addresses. The control plane reaches CDP through an SSH
  # session to the browser pod (port 22 is the only externally exposed port);
  # ssh.Client.Dial then opens 127.0.0.1:9222 inside the pod.
  # Ensure Downloads directory exists. On the on-demand layout this path is
  # the agent's home PVC mounted via subPath so files Chromium saves here are
  # immediately visible to OpenClaw and the agent terminal.
  s6-setuidgid claworc mkdir -p /home/claworc/Downloads "$PROFILE_DIR" 2>/dev/null || true
  s6-setuidgid claworc $BROWSER \
    --no-first-run \
    --password-store=basic \
    --simulate-outdated-no-au='Tue, 31 Dec 2099 23:59:59 GMT' \
    --window-size="${DISPLAY_WIDTH:-1920},${DISPLAY_HEIGHT:-1080}" \
    --window-position=0,0 \
    --user-data-dir=/home/claworc/chrome-data \
    --download-directory=/home/claworc/Downloads \
    --remote-debugging-port=9222 \
    --remote-debugging-address=127.0.0.1 \
    --disable-default-apps \
    --disable-features=CloseWindowWithLastTab \
    --disable-blink-features=AutomationControlled \
    --disable-infobars \
    --disable-component-update \
    --load-extension=/opt/stealth-extension \
    "${EXTRA_ARGS[@]}" > /dev/null 2>&1
  sleep 1
done
